# systemctl.service
[Unit]
Description=CommitPush
After=network.target

[Service]
Type=forking
ExecStart=/pathtobin